# system_services/rover-base-station-gstreamer.timer
[Timer]
OnBootSec=20s

[Install]
WantedBy=multi-user.target